# useful/profile.d/xilinx-opencl.sh
# User specific aliases and functions

source /opt/rh/rh-git29/enable

export LM_LICENSE_FILE=/opt/Xilinx/xilinx.lic

#source /opt/Xilinx/SDx/2017.1/settings64.sh

export GPUDMA_DIR=/xprj/gpudma

source /opt/rh/devtoolset-6/enable 

source /opt/Xilinx/ku3_dsa_2017.1-xpr/xbinst/setup.sh



